# units/slp/var-log.automount
[Unit]
Description=/var/log Automount Point

[Automount]
Where=/var/log
