# units/slp/var-lib.mount
[Unit]
Description=/var/lib Mount Point
Before=local-fs.target

[Mount]
What=/opt/var/lib
Where=/var/lib
Type=bind
Options=bind
