# units/slp/var-log.mount
[Unit]
Description=/var/log Mount Point
Before=local-fs.target

[Mount]
What=/opt/var/log
Where=/var/log
Type=bind
Options=bind
