# units/slp/var-lib.automount
[Unit]
Description=/var/lib Automount Point

[Automount]
Where=/var/lib
