# units/slp/var-cache.automount
[Unit]
Description=/var/cache Automount Point

[Automount]
Where=/var/cache
