# units/slp/var-lock.mount
[Unit]
Description=/var/lock Mount Point
Before=local-fs.target

[Mount]
What=/run/lock
Where=/var/lock
Type=bind
Options=bind
